# debian/galmon.galmon-upgrade.timer
[Unit]
Description=Update galmon
Documentation=man:ubxtool

[Timer]
OnBootSec=3min
OnUnitActiveSec=3d

[Install]
WantedBy=timers.target
